# services/kde-connectd.service
[Unit]
Description=KDE Connect Daemon
After=network.target

[Service]
ExecStart=/usr/bin/kdeconnectd
Restart=on-failure

[Install]
WantedBy=default.target
